# crontab
#
# 0 0 * * 0,4 cron_start.sh > /proc/1/fd/1 2>&1
0 0 * * 0,4 /usr/local/bin/python3 /app/main.py > /proc/1/fd/1 2>&1
#